# ansible/roles/stash-face-tagger/templates/face-tagger.timer.j2
[Unit]
Description=Stash Face Tagger — auto-tag scenes with performer faces

[Timer]
OnCalendar=*-*-* 02:00:00
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
